# crontab
0,30 * * * * python ~/../datalogger/get_weather.py >> ~/../datalogger/log/crontab_log.txt